# docs/supervisord.conf
[program:boat]
directory=/www/boatcontrol
command=node index.js
autostart=true
autorestart=true
stdout_logfile=/www/boatcontrol/boat.log
stderr_logfile=/www/boatcontrol/error.log
